# dpkg/phalanx.v1.0/phalanx-installation_0.1_amd64/etc/systemd/system/mrpc.service
[Unit]
Description=msf rpc service

[Service]
Type=forking
ExecStart=/usr/bin/msfrpcd -U user -P pass

[Install]
WantedBy=multi-user.target