# dpkg/phalanx.v1.0/phalanx-installation_0.1_amd64/etc/systemd/system/sqlmap_server.service
[Unit]
Description=Sqlmap Server Service
After=jobmanager.service

[Service]
Type=simple
WorkingDirectory=/opt/crawler/
ExecStart=/usr/bin/python3 /opt/crawler/server.py
Restart=on-failure
StandardOutput=file:/var/log/phalanx/sqlmap_server.log
StandardError=file:/var/log/phalanx/sqlmap_server.log

[Install]
WantedBy=multi-user.target
