# dpkg/phalanx.v1.0/phalanx-installation_0.1_amd64/etc/systemd/system/tunnel.service
[Unit]
Description=Tunnel Service
After=jobmanager.service

[Service]
Type=simple
WorkingDirectory=/opt/tunnel/
ExecStart=/opt/tunnel/tunnel-service
Restart=on-failure
StandardOutput=file:/var/log/phalanx/tunnel.log
StandardError=file:/var/log/phalanx/tunnel.log

[Install]
WantedBy=multi-user.target
